# build/docker/supervisor/ini/sbot.ini
[program:sbot]
command=/root/gewe/api/sbot start --config=/root/gewe/api/config/robot.yaml
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/root/gewe/api/log/sbot.log
stopwaitsecs=3600

