# docker/supervisor/ini/gog7a6v90.ini
[program:gog7a6v90]
directory=/app
command=/app/gog7a6v90
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
startsecs=10
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=20
stdout_logfile=/app/gog7a6v90.log