# docker/supervisor/ini/khan.ini
[program:khan]
command=/app/khan start --config=/app/config/robot.yaml
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/app/khan.log
stopwaitsecs=3600

