# docker/supervisor/ini/sbot.ini
[program:sbot]
command=/app/sbot start --config=/app/config/robot.yaml
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/app/sbot.log
stopwaitsecs=3600

